# systems/kafka/1.1/build-image/supervisord.conf
[supervisord]
nodaemon=true
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D

[program:kafka]
command=/bin/bash -c "/usr/local/bin/kafka-clusternode"
