# systems/zk/branch-3.6/build-image/supervisord.conf
[supervisord]
nodaemon=true

[program:zookeeper]
command=/bin/bash -c "/usr/local/bin/zk-clusternode"
