# systems/hdfs/release-2.9.0/build-image/supervisord.conf
[supervisord]
nodaemon=true
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D

[program:hadoop]
command=/bin/bash -c "/usr/local/bin/hadoop-startup"
