# systems/cass/cassandra-2.1.0/build-image/supervisord.conf
[supervisord]
nodaemon=true
autorestart=false

[program:cassandra]
command=/bin/bash -c "/usr/local/bin/cassandra-clusternode"
