# systems/hive/template/build-image/supervisord.conf
[supervisord]
nodaemon=true
autorestart=false

[program:sshd]
command=/usr/sbin/sshd -D

[program:hive]
command=/bin/bash -c "/usr/local/bin/hive-startup"
